# cron/crontab
PATH=/usr/local/bin:/usr/bin:/bin
* * * * * env DISPLAY=:0 /home/unkiwii/dotfiles/cron/alert-on-low-battery
* * * * * env DISPLAY=:0 /home/unkiwii/dotfiles/cron/alert-on-high-disk-usage
